# packages/elem-hasplib/elem_hasplib-1.0.0-py3-none-any.whl/elem_hasplib/aksusbd/bin/aksusbd.service
[Unit]
Description=Sentinel LDK Runtime Environment (aksusbd daemon)
Before=hasplmd.service
Requires=hasplmd.service

[Service]
Type=forking
ExecStart=/usr/sbin/aksusbd

[Install]
WantedBy=multi-user.target

