# packages/elem-hasplib/elem_hasplib-1.0.0-py3-none-any.whl/elem_hasplib/aksusbd/bin/hasplmd_arm64.service
[Unit]
Description=Sentinel LDK Runtime Environment (hasplmd daemon)

[Service]
Type=forking
ExecStart=/usr/sbin/hasplmd_arm64 -s

[Install]
WantedBy=multi-user.target
