# packages/elem-hasplib/elem_hasplib-1.0.0-py3-none-any.whl/elem_hasplib/aksusbd/bin/hasplmd.service
[Unit]
Description=Sentinel LDK Runtime Environment (hasplmd daemon)
Requires=aksusbd.service
After=aksusbd.service

[Service]
Type=forking
ExecStart=/usr/sbin/hasplmd -s

[Install]
WantedBy=multi-user.target
